# gunicorn/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn -k gevent --workers 4 -b localhost:5000 run:app
directory=/var/www/app